# supervisord.production.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=nginx -g "daemon off;"
priority=10
stdout_logfile=/app/logs/nginx_stdout.log
stderr_logfile=/app/logs/nginx_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
autorestart=true
startretries=3

[program:rust-backend]
command=/app/target/release/visionclaw-server --port 4001 --gpu-debug
priority=20
directory=/app
environment=RUST_LOG="%(ENV_RUST_LOG)s",RUST_BACKTRACE="1",NVIDIA_VISIBLE_DEVICES="%(ENV_NVIDIA_VISIBLE_DEVICES)s",NVIDIA_DRIVER_CAPABILITIES="compute,utility"
stdout_logfile=/app/logs/rust_stdout.log
stderr_logfile=/app/logs/rust_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
autorestart=true
startretries=5
startsecs=10
stopasgroup=true
killasgroup=true