# multi-agent-docker/unified-config/supervisord.unified.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================================================
# PostgreSQL Database (Port 5432) - RuVector Vector Storage
# DISABLED: Using external ruvector-postgres on docker_ragflow network instead
# Start manually if needed: supervisorctl start postgresql
# ============================================================================

[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgres/data
user=postgres
environment=PGDATA="/var/lib/postgres/data"
autostart=false
autorestart=true
priority=50
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql.error.log
stopasgroup=true
killasgroup=true

# ============================================================================
# Xvfb + x11vnc with VirtualGL (Port 5901) - 4K Resolution with GPU
# ============================================================================

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 2048x2048x24 -ac +extension GLX +render -noreset
user=root
environment=HOME="/root",DISPLAY=":1"
autostart=true
autorestart=true
priority=90
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.error.log
stopasgroup=true
killasgroup=true

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -rfbport 5901 -forever -shared -nopw -xkb -cursor arrow -cursor_drag -repeat -nowcr
user=root
environment=HOME="/root",DISPLAY=":1"
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.error.log
depends_on=xvfb
stopasgroup=true
killasgroup=true

# ============================================================================
# XFCE4 Desktop Environment (4K with Terminal Grid Autostart)
# ============================================================================

[program:openbox]
command=/usr/sbin/openbox
user=devuser
environment=HOME="/home/devuser",USER="devuser",DISPLAY=":1"
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/openbox.log
stderr_logfile=/var/log/openbox.error.log
depends_on=x11vnc
startsecs=5
stopasgroup=true
killasgroup=true

[program:tint2]
command=/usr/sbin/tint2
user=devuser
environment=HOME="/home/devuser",USER="devuser",DISPLAY=":1"
autostart=true
autorestart=true
priority=210
stdout_logfile=/var/log/tint2.log
stderr_logfile=/var/log/tint2.error.log
depends_on=openbox
startsecs=3
stopasgroup=true
killasgroup=true

# ============================================================================
# Terminal Grid Autostart (7 kitty terminals in grid layout)
# ============================================================================

[program:terminal-grid]
command=/home/devuser/.config/autostart-terminals.sh
user=devuser
environment=HOME="/home/devuser",USER="devuser",DISPLAY=":1",DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/1000/bus",LIBGL_ALWAYS_SOFTWARE="1",GALLIUM_DRIVER="llvmpipe"
autostart=true
autorestart=false
priority=300
stdout_logfile=/var/log/terminal-grid.log
stderr_logfile=/var/log/terminal-grid.error.log
depends_on=openbox
startsecs=0

# ============================================================================
# SSH Server (Port 22)
# ============================================================================

[program:sshd]
command=/usr/sbin/sshd -D -e
user=root
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
priority=50
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.error.log

# ============================================================================
# Management API (Port 9090)
# ============================================================================

[program:management-api]
command=/usr/local/bin/node /opt/management-api/server.js
directory=/opt/management-api
user=devuser
environment=HOME="/home/devuser",PORT="9090",NODE_ENV="production"
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
priority=300
stdout_logfile=/var/log/management-api.log
stderr_logfile=/var/log/management-api.error.log

# ============================================================================
# code-server (Port 8080)
# ============================================================================

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /home/devuser/workspace
user=devuser
environment=HOME="/home/devuser",USER="devuser"
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
priority=400
stdout_logfile=/var/log/code-server.log
stderr_logfile=/var/log/code-server.error.log

# ============================================================================
# Claude Z.AI Service (Port 9600 - Internal Only)
# Runs as zai-user for isolation
# NOTE: Set ZAI_EXTERNAL=true in environment to disable (when using external container)
# ============================================================================

[program:claude-zai]
command=/usr/local/bin/node /opt/claude-zai/server.js
directory=/opt/claude-zai
user=zai-user
environment=HOME="/home/zai-user",PORT="9600",NODE_ENV="production"
autostart=%(ENV_ZAI_INTERNAL)s
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
priority=500
stdout_logfile=/var/log/claude-zai.log
stderr_logfile=/var/log/claude-zai.error.log

# ============================================================================
# MCP Gateway (TCP + WebSocket) - VisionFlow Integration
# Provides MCP TCP on port 9500 and WebSocket on port 3002
# ============================================================================

[program:mcp-gateway]
command=/usr/local/bin/node /opt/mcp-gateway/server.js
directory=/opt/mcp-gateway
user=devuser
environment=HOME="/home/devuser",MCP_TCP_PORT="9500",MCP_WS_PORT="3002",LOG_LEVEL="info"
autostart=true
autorestart=true
priority=150
stdout_logfile=/var/log/mcp-gateway.log
stderr_logfile=/var/log/mcp-gateway.error.log

# ============================================================================
# MCP Servers (devuser) - Skill Integration
# ============================================================================

# Web Summary MCP - DEMOTED to on-demand (autostart=false)
# Start manually when web summarization needed: supervisorctl start web-summary-mcp
[program:web-summary-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/web-summary/mcp-server/server.js
directory=/home/devuser/.claude/skills/web-summary/mcp-server
user=devuser
environment=HOME="/home/devuser",WEB_SUMMARY_TOOL_PATH="/home/devuser/.claude/skills/web-summary/tools/web_summary_tool.py",ZAI_CONTAINER_URL="http://localhost:9600"
autostart=false
autorestart=true
priority=510
stdout_logfile=/var/log/web-summary-mcp.log
stderr_logfile=/var/log/web-summary-mcp.error.log

[program:qgis-mcp]
command=/usr/bin/python3 -u /home/devuser/.claude/skills/qgis/tools/qgis_mcp.py
directory=/home/devuser/.claude/skills/qgis/tools
user=devuser
environment=HOME="/home/devuser",QGIS_HOST="localhost",QGIS_PORT="9877"
autostart=true
autorestart=true
priority=511
stdout_logfile=/var/log/qgis-mcp.log
stderr_logfile=/var/log/qgis-mcp.error.log

[program:blender-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/blender/tools/mcp-blender-client.js
directory=/home/devuser/.claude/skills/blender/tools
user=devuser
environment=HOME="/home/devuser",BLENDER_HOST="localhost",BLENDER_PORT="9876"
autostart=true
autorestart=true
priority=512
stdout_logfile=/var/log/blender-mcp.log
stderr_logfile=/var/log/blender-mcp.error.log

# ImageMagick MCP - DEMOTED to on-demand (autostart=false)
# Start manually when image processing needed: supervisorctl start imagemagick-mcp
[program:imagemagick-mcp]
command=/usr/bin/python3 -u /home/devuser/.claude/skills/imagemagick/tools/imagemagick_mcp.py
directory=/home/devuser/.claude/skills/imagemagick/tools
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=513
stdout_logfile=/var/log/imagemagick-mcp.log
stderr_logfile=/var/log/imagemagick-mcp.error.log

# ============================================================================
# ComfyUI - REMOVED from internal deployment
# Now runs as external Docker container, accessed via comfyui skill MCP
# ============================================================================

[program:kicad-mcp]
command=/usr/bin/python3 -u /home/devuser/.claude/skills/kicad/tools/kicad_mcp.py
directory=/home/devuser/.claude/skills/kicad/tools
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=514
stdout_logfile=/var/log/kicad-mcp.log
stderr_logfile=/var/log/kicad-mcp.error.log

[program:ngspice-mcp]
command=/usr/bin/python3 -u /home/devuser/.claude/skills/ngspice/tools/ngspice_mcp.py
directory=/home/devuser/.claude/skills/ngspice/tools
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=515
stdout_logfile=/var/log/ngspice-mcp.log
stderr_logfile=/var/log/ngspice-mcp.error.log

[program:pbr-mcp]
command=/usr/bin/python3 -u /home/devuser/.claude/skills/pbr-rendering/tools/pbr_mcp_client.py
directory=/home/devuser/.claude/skills/pbr-rendering/tools
user=devuser
environment=HOME="/home/devuser",PBR_HOST="localhost",PBR_PORT="9878"
autostart=false
autorestart=true
priority=516
stdout_logfile=/var/log/pbr-mcp.log
stderr_logfile=/var/log/pbr-mcp.error.log

# ============================================================================
# @claude-flow/browser - PRIMARY Browser Automation (V3)
# DISABLED: Package has no standalone 'serve' command - browser available via
# ruflo MCP server's browser_* tools instead (59 tools built-in)
# ============================================================================

[program:ruflo-browser]
command=/usr/local/bin/npx @claude-flow/browser serve --port 9510
directory=/home/devuser
user=devuser
environment=HOME="/home/devuser",DISPLAY=":1",BROWSER_HEADLESS="false"
autostart=false
autorestart=true
priority=515
stdout_logfile=/var/log/claude-flow-browser.log
stderr_logfile=/var/log/claude-flow-browser.error.log

# Playwright MCP - DEMOTED to fallback (autostart=false)
# Use @claude-flow/browser as primary browser system
[program:playwright-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/playwright/tools/playwright-mcp-local.js
directory=/home/devuser/.claude/skills/playwright/tools
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=517
stdout_logfile=/var/log/playwright-mcp.log
stderr_logfile=/var/log/playwright-mcp.error.log

[program:perplexity-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/perplexity/mcp-server/server.js
directory=/home/devuser/.claude/skills/perplexity/mcp-server
user=devuser
environment=HOME="/home/devuser"
autostart=false
autorestart=true
priority=518
stdout_logfile=/var/log/perplexity-mcp.log
stderr_logfile=/var/log/perplexity-mcp.error.log

# Jupyter Notebooks MCP - DEMOTED to on-demand (autostart=false)
# Start manually when notebooks needed: supervisorctl start jupyter-notebooks-mcp
[program:jupyter-notebooks-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/jupyter-notebooks/server.js
directory=/home/devuser/.claude/skills/jupyter-notebooks
user=devuser
environment=HOME="/home/devuser",JUPYTER_NOTEBOOK_DIR="/home/devuser/workspace"
autostart=false
autorestart=true
priority=520
stdout_logfile=/var/log/jupyter-notebooks-mcp.log
stderr_logfile=/var/log/jupyter-notebooks-mcp.error.log

# ============================================================================
# OpenAI Codex MCP Server (GPT-5.4)
# Exposes codex_generate and codex_review tools via MCP for cross-agent use
# Runs as devuser (skill files live under ~/.claude/skills which is devuser-owned)
# API isolation via OPENAI_API_KEY env var, not unix user
# ============================================================================

[program:openai-codex-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/openai-codex/mcp-server/server.js
directory=/home/devuser/.claude/skills/openai-codex/mcp-server
user=devuser
environment=HOME="/home/devuser",USER="devuser",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",OPENAI_DEFAULT_MODEL="%(ENV_OPENAI_DEFAULT_MODEL)s"
autostart=true
autorestart=true
priority=605
stdout_logfile=/var/log/openai-codex-mcp.log
stderr_logfile=/var/log/openai-codex-mcp.error.log

# ============================================================================
# Gemini Flow Daemon (gemini-user)
# ============================================================================

[program:gemini-flow]
command=/usr/local/bin/gemini-flow mcp start
directory=/home/gemini-user/workspace
user=gemini-user
environment=HOME="/home/gemini-user",USER="gemini-user",MCP_SOCKET="/var/run/agentic-services/gemini-mcp.sock",GEMINI_FLOW_DATA_DIR="/home/gemini-user/.gemini-flow"
autostart=true
autorestart=true
priority=600
stdout_logfile=/var/log/gemini-flow.log
stderr_logfile=/var/log/gemini-flow.error.log

# ============================================================================
# Screensaver Disable (runs after XFCE4 starts)
# ============================================================================

[program:disable-screensaver]
command=/home/devuser/.config/disable-screensaver.sh
user=devuser
environment=HOME="/home/devuser",USER="devuser",DISPLAY=":1"
autostart=true
autorestart=false
priority=250
stdout_logfile=/var/log/disable-screensaver.log
stderr_logfile=/var/log/disable-screensaver.error.log
depends_on=openbox
startsecs=0

# ============================================================================
# tmux Workspace Auto-Start (devuser)
# Starts after a delay to ensure other services are ready
# ============================================================================

[program:tmux-autostart]
command=/bin/bash -c "sleep 15 && /home/devuser/.config/tmux-autostart.sh"
user=devuser
environment=HOME="/home/devuser",USER="devuser",WORKSPACE="/home/devuser/workspace",AGENTS_DIR="/home/devuser/agents"
autostart=true
autorestart=false
priority=900
stdout_logfile=/var/log/tmux-autostart.log
stderr_logfile=/var/log/tmux-autostart.error.log
depends_on=openbox

# ============================================================================
# Management API Health Check (one-time verification)
# Runs 10 seconds after startup to verify Management API is responding
# ============================================================================

[program:management-api-healthcheck]
command=/bin/bash /opt/scripts/verify-management-api.sh
user=root
environment=MANAGEMENT_API_HOST="localhost",MANAGEMENT_API_PORT="9090"
autostart=true
autorestart=false
startsecs=0
priority=950
stdout_logfile=/var/log/management-api-healthcheck.log
stderr_logfile=/var/log/management-api-healthcheck.error.log

# ============================================================================
# DBus Daemon (required for desktop apps)
# ============================================================================

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork
user=root
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/dbus.log
stderr_logfile=/var/log/dbus.error.log

[program:dbus-user]
command=/usr/bin/dbus-daemon --session --nofork --address=unix:path=/run/user/1000/bus
user=devuser
environment=HOME="/home/devuser",USER="devuser",XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/dbus-user.log
stderr_logfile=/var/log/dbus-user.error.log

# ============================================================================
# HTTPS Bridge Proxy (Bridges https://localhost:3001 -> host http:3001)
# Enables secure browser connections to host development servers
# ============================================================================

[program:https-bridge]
command=/usr/local/bin/node /opt/https-bridge/https-proxy.js
directory=/opt/https-bridge
user=devuser
environment=HOME="/home/devuser",HOST_IP="host.docker.internal",HTTPS_PORT="3001",TARGET_PORT="3001",CERT_DIR="/opt/https-bridge"
autostart=true
autorestart=true
priority=350
stdout_logfile=/var/log/https-bridge.log
stderr_logfile=/var/log/https-bridge.error.log

# ============================================================================
# Local LLM Proxy - Anthropic-to-OpenAI translation for Nemotron 3 120B
# Allows Claude CLI to communicate with llama.cpp via ANTHROPIC_BASE_URL
# ============================================================================

[program:local-llm-proxy]
command=/usr/local/bin/node /opt/scripts/local-llm-proxy.mjs
directory=/home/local-private
user=local-private
environment=HOME="/home/local-private",LOCAL_LLM_HOST="%(ENV_LOCAL_LLM_HOST)s",LOCAL_LLM_PORT="%(ENV_LOCAL_LLM_PORT)s",LOCAL_LLM_MODEL="%(ENV_LOCAL_LLM_MODEL)s",PROXY_PORT="3100",NODE_ENV="production",NODE_PATH="/usr/local/lib/node_modules"
autostart=false
autorestart=true
priority=360
stdout_logfile=/var/log/local-llm-proxy.log
stderr_logfile=/var/log/local-llm-proxy.error.log

# ============================================================================
# Claude Telegram Mirror (ctm) — Bidirectional Claude Code <-> Telegram
# Daemon bridges hook events to Telegram and relays Telegram replies to CLI.
# Set TELEGRAM_MIRROR=true in .env + provide TELEGRAM_BOT_TOKEN & TELEGRAM_CHAT_ID
# ============================================================================

[program:claude-telegram-mirror]
command=/usr/local/bin/ctm start
directory=/home/devuser
user=devuser
environment=HOME="/home/devuser",RUST_LOG="info",TELEGRAM_MIRROR="%(ENV_TELEGRAM_MIRROR)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s",CTM_LLM_SUMMARIZE_URL="%(ENV_CTM_LLM_SUMMARIZE_URL)s",CTM_LLM_API_KEY="%(ENV_CTM_LLM_API_KEY)s"
autostart=%(ENV_TELEGRAM_MIRROR)s
autorestart=true
priority=550
stdout_logfile=/var/log/claude-telegram-mirror.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/claude-telegram-mirror.error.log
stderr_logfile_maxbytes=5MB

[program:host-webserver-debug-mcp]
command=/usr/local/bin/node /home/devuser/.claude/skills/host-webserver-debug/mcp-server/server.js
directory=/home/devuser/.claude/skills/host-webserver-debug/mcp-server
user=devuser
environment=HOME="/home/devuser",DISPLAY=":1",HOST_GATEWAY_IP="host.docker.internal",SCREENSHOT_DIR="/tmp/screenshots"
autostart=false
autorestart=true
priority=519
stdout_logfile=/var/log/host-webserver-debug-mcp.log
stderr_logfile=/var/log/host-webserver-debug-mcp.error.log
